# opennic-up.timer
[Unit]
Description=Schedule opennic-up DNS update

[Timer]
OnCalendar=Wed,Sat *-*-* 11:00:00
Persistent=true

[Install]
WantedBy=timers.target
